# zones/mapi/fs.root/site/postgresql.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM 
"/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!-- 
Copyright 2006 Sun Microsystems, Inc. All rights reserved. 
Use is subject to license terms. 

CDDL HEADER START 

The contents of this file are subject to the terms of the 
Common Development and Distribution License (the "License"). 
You may not use this file except in compliance with the License. 

You can obtain a copy of the license at usr/src/OPENSOLARIS.LICENSE 
or http://www.opensolaris.org/os/licensing. 
See the License for the specific language governing permissions 
and limitations under the License. 

When distributing Covered Code, include this CDDL HEADER in each
file and include the License file at usr/src/OPENSOLARIS.LICENSE. 

If applicable, add the following below this CDDL HEADER, with the 
fields enclosed by brackets "[]" replaced with your own identifying 
information: Portions Copyright [yyyy] [name of copyright owner] 

CDDL HEADER END 

     ident "%Z%%M% %I% %E SMI" 

     NOTE: This service manifest is not editable; its contents will 
     be overwritten by package or patch 
     operations, including operating system upgrade. Make customizations 
     in a different file. 
--> 

<service_bundle type='manifest' name='postgresql'>

<service 
     name='network/postgresql' 
     type='service' 
     version='1'>

     <!-- 
     Wait for network interfaces to be initialized. 
     --> 
     <dependency 
          name='network' 
          grouping='require_all' 
          restart_on='none' 
          type='service'> 
          <service_fmri value='svc:/milestone/network:default' /> 
     </dependency> 

     <!-- 
          Wait for all local filesystems to be mounted. 
     --> 
     <dependency 
          name='filesystem-local' 
          grouping='require_all' 
          restart_on='none' 
          type='service'> 
     <service_fmri value='svc:/system/filesystem/local:default' /> 
     </dependency>

     <exec_method 
          type='method' 
          name='start' 
          exec='/opt/local/svc/svc-postgresql start' 
          timeout_seconds='300' />
     <exec_method 
          type='method' 
          name='stop' 
          exec='/opt/local/svc/svc-postgresql stop' 
          timeout_seconds='300' /> 

     <exec_method 
          type='method' 
          name='refresh' 
          exec='/opt/local/svc/svc-postgresql refresh' 
          timeout_seconds='60' /> 

     <!-- 
       We define two instances of PostgreSQL as examples. 
     --> 

     <instance name='default' enabled='true'> 
          <method_context> 
               <method_credential user='postgres' group='postgres' /> 
          </method_context>

          <!-- 
             Make sure the data configurable property points to the 
             appropriate database directory. 
          --> 

          <property_group name='postgresql' type='application'> 
               <propval name='data' type='astring' 
                  value='/var/pgsql/data' /> 
               <propval name='log' type='astring' 
                  value='/var/log/postgresql.log' /> 
          </property_group> 

     </instance> 

     <stability value='Evolving' /> 

     <template> 
          <common_name> 
               <loctext xml:lang='C'> 
               PostgreSQL RDBMS 
               </loctext> 
          </common_name> 
          <documentation> 
               <manpage title='postgres' section='1M' /> 
               <doc_link name='postgresql.org' 
                    uri='http://postgresql.org' /> 
          </documentation> 
     </template> 

</service> 

</service_bundle>
